# hardware/device-hub/systemd/4dx-hub.service
# systemd service file for 4DX@HOME Device Hub
# TODO: Configure for production deployment

[Unit]
Description=4DX@HOME Device Hub Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/4dx-home/hardware/device-hub
ExecStart=/usr/bin/python3 /home/pi/4dx-home/hardware/device-hub/src/hub.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/home/pi/4dx-home/hardware/device-hub

[Install]
WantedBy=multi-user.target